# nbackup/nbackup_install/backup/backup.service
[Unit]
Description=WMS_Backup

[Service]
Type=simple
ExecStart=/data/backup/run.sh

[Install]
WantedBy=multi-user.target
